# mem-watchdog.service
[Unit]
Description=Crostini Memory Watchdog — protects VS Code from OOM killer
Documentation=https://github.com/chf3198/crostini-mem-watchdog
After=default.target

# Why this service exists:
# earlyoom v1.7 crashes on Crostini immediately (exit 104) because the
# ChromeOS kernel reports SwapFree: 18446744073709551360 kB (uint64 overflow).
# earlyoom's strtol() call overflows fatally. This custom watchdog uses only
# MemAvailable and PSI, which are both correct on this kernel.

[Service]
Type=simple
ExecStart=%h/.local/bin/mem-watchdog.sh
Restart=always
RestartSec=5
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal
# Rate limiting: DISABLED at the per-unit level (LogRateLimitBurst=0).
# The daemon is architecturally self-rate-limiting — it logs only on action
# events (SIGTERM/SIGKILL, threshold crossings, startup mode activation), never
# on every 2s/0.5s poll cycle. A per-unit burst cap would risk silently dropping
# CRIT-priority kill-action log lines during an OOM cascade — the exact moment
# log integrity matters most. The system-wide cap in journald-limits.conf
# (RateLimitBurst=500/30s) is sufficient protection against runaway services.
LogRateLimitBurst=0

[Install]
WantedBy=default.target
