# mem-watchdog.service
[Unit]
Description=Crostini Memory Watchdog — protects VS Code from OOM killer
Documentation=https://github.com/chf3198/crostini-mem-watchdog
After=default.target

# Why this service exists:
# earlyoom v1.7 crashes on Crostini immediately (exit 104) because the
# ChromeOS kernel reports SwapFree: 18446744073709551360 kB (uint64 overflow).
# earlyoom's strtol() call overflows fatally. This custom watchdog uses only
# MemAvailable and PSI, which are both correct on this kernel.

[Service]
Type=simple
ExecStart=%h/.local/bin/mem-watchdog.sh
Restart=always
RestartSec=5
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal
# Rate-limit journald writes from this service.
# Allows 100 log entries per 30-second window — sufficient for normal operation
# (daemon logs only on action events, not every poll cycle) while preventing
# a flood during startup mode (0.5 s polling × 90 s = up to 180 cycles where
# OOM score adjustments for new PIDs can each emit a log line).
LogRateLimitIntervalSec=30s
LogRateLimitBurst=100

[Install]
WantedBy=default.target
